# examples/opensensor-space/systemd/opensensor_space_systemd.service
[Unit]
Description=OpenSensor Space Data Collection
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=REPLACE_USER
Group=REPLACE_USER
WorkingDirectory=REPLACE_PROJECT_ROOT
Environment="PATH=REPLACE_PROJECT_ROOT/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Main data collection service
ExecStart=REPLACE_PROJECT_ROOT/.venv/bin/python examples/opensensor-space/opensensor_space.py

# Restart policy
Restart=always
RestartSec=60

# Logging
StandardOutput=append:REPLACE_PROJECT_ROOT/logs/opensensor_space.log
StandardError=append:REPLACE_PROJECT_ROOT/logs/opensensor_space_error.log

# Security hardening (optional - uncomment if needed)
# NoNewPrivileges=true
# PrivateTmp=true

[Install]
WantedBy=multi-user.target
