# examples/opensensor-space/systemd/sync_timer.service
[Unit]
Description=OpenSensor Space Storage Sync
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=REPLACE_USER
Group=REPLACE_USER
WorkingDirectory=REPLACE_PROJECT_ROOT
ExecStart=REPLACE_PROJECT_ROOT/examples/opensensor-space/sync_to_storage.sh

# Logging
StandardOutput=append:REPLACE_PROJECT_ROOT/logs/sync_service.log
StandardError=append:REPLACE_PROJECT_ROOT/logs/sync_service_error.log
